# deploy/systemd/krisnabot.service.example
[Unit]
Description=KRISNABOT FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/krisnabot
EnvironmentFile=/opt/krisnabot/.env
ExecStart=/opt/krisnabot/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
